# backend/deploy/supervisor.conf
[program:petglow-api]
directory=/srv/petglow/backend
command=/srv/petglow/backend/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
autostart=true
autorestart=true
stderr_logfile=/var/log/petglow-api.err.log
stdout_logfile=/var/log/petglow-api.out.log
stopasgroup=true
killasgroup=true
